# dnd-toolkit.service
[Unit]
Description=DnD DM Toolkit
After=network.target

[Service]
Type=simple
User=bgmaddox
WorkingDirectory=/home/bgmaddox/dnd
EnvironmentFile=/home/bgmaddox/dnd/.env
ExecStart=/home/bgmaddox/dnd/.venv/bin/uvicorn server:app --host 0.0.0.0 --port 8502
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
